# dotfiles/systemd/.config/systemd/user/remote-plasma-gui.target
[Unit]
Description=Remote Plasma GUI Session

Wants=xdummy.service
Wants=plasma-remote.service
